# psplash-start.service
[Unit]
Description=Starts Psplash Boot screen
DefaultDependencies=no
RequiresMountsFor=/run
Before=sysinit.target
After=dev-fb0.device
Requires=dev-fb0.device

[Service]
Type=exec
ExecStartPre=/bin/sh -c 'while [ ! -e /dev/fb0 ]; do sleep 0.1; done'
ExecStart=/usr/local/bin/psplash -n --fbdev 0

[Install]
WantedBy=sysinit.target
